# rootfs/etc/services.d/yandex-disk/run
#!/usr/bin/with-contenv bash

yandex-disk start --dir $YANDEX_FOLDER -D