# rootfs/etc/services.d/reflex/run
#!/usr/bin/with-contenv bash

cd $FLEXGET_SYNC_DIR
/usr/sbin/reflex -g '*.*' /usr/sbin/sync.sh